# chap4/c4supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/bin/bash -c "echo 'start sshd'"
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true